# roles/home_server/templates/librespeed.service.j2
[Unit]
Description=LibreSpeed bandwidth test (writes node-exporter textfile metrics)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart={{ librespeed_script }}
Nice=10
